# scripts/launchd/com.gosenderr.cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gosenderr.cleanup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/<your-username>/.local/bin/gosenderr-cleanup.sh</string>
    <string>--yes</string>
    <string>--pnpm-prune</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
    <key>Weekday</key>
    <integer>1</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/<your-username>/Library/Logs/gosenderr-cleanup.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/<your-username>/Library/Logs/gosenderr-cleanup.err</string>
</dict>
</plist>